# .config/autostart/setxkb.desktop
# See https://gitlab.gnome.org/GNOME/gnome-shell/-/issues/1154
[Desktop Entry]
Type=Application
Exec=setxkb us,ru
